# backend/docker/supervisord.conf
# Configuration Supervisor : permet de lancer et maintenir plusieurs services dans le même conteneur Docker
# ici PHP-FPM (exécution PHP) et Nginx (serveur web), avec redémarrage automatique en cas de crash.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; Exécute les commandes d'initialisation Laravel (cache config + migrations)
[program:laravel-init]
command=sh -c "php artisan config:cache && php artisan migrate --force"
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Démarre PHP-FPM
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Démarre nginx
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0